# config/profile/profile.d/z-10-podman.sh
if [ -f /etc/containers/storage.conf ]; then
    export CONTAINERS_STORAGE_CONF=/etc/containers/storage.conf
fi
